# php-start-4000.service
[Unit]

Description= Starts PHP for External-Apps-Launcher

[Service]

ExecStart=/bin/bash /usr/bin/test.sh


[Install]
WantedBy=multi-user.target
