# php-start.service
[Unit]

Description= Starts PHP for Smith-Pad-OS-Shell

[Service]

ExecStart=/bin/bash /usr/bin/php-start.sh


[Install]
WantedBy=multi-user.target
