# template.service
[Unit]

Description=template

[Service]

ExecStart=/bin/bash /usr/sbin


[Install]
WantedBy=multi-user.target
