# Smith-Pad-OS-Shell-Configs/test.service
[Unit]

Description= Description Here

[Service]

ExecStart=/bin/bash /usr/bin/test.sh


[Install]
WantedBy=multi-user.target
