# etc/systemd/template.service
[Unit]

Description=template

[Service]

ExecStart=/bin/bash /usr/bin


[Install]
WantedBy=multi-user.target
